# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:mariadb]
command=/usr/sbin/mysqld --user=mysql
autorestart=true
priority=5

[program:php-fpm]
command=/bin/sh -c "until mysqladmin ping -h 127.0.0.1 --silent 2>/dev/null; do sleep 1; done; exec /usr/local/sbin/php-fpm"
autorestart=true
priority=10

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
priority=20
